# service-conf/staging/gunicorn-cyph3r.socket
Description=Cyph3r Gunicorn Socket
[Socket]
ListenStream=/run/gunicorn-cyph3r.sock
User=nginx
Mode=600
[Install]
WantedBy=sockets.target